# deploy/roles/filter/templates/supervisor/osp-network.conf.j2

[program:osp-network]
command={{ filter_src }}/env/bin/uwsgi --ini {{ uwsgi_config }}/osp-network.ini
autostart={{ supervisor_autostart }}
