# deploy/roles/network/templates/supervisor/osp-network.conf.j2

[program:osp-network]
command={{ network_src }}/env/bin/uwsgi --ini {{ uwsgi_config }}/osp-network.ini
autostart={{ supervisor_autostart }}
